# limine/spectrumos-limine-sync.path
[Unit]
Description=Watch for SpectrumOS wallpaper changes
After=local-fs.target

[Path]
PathChanged=/var/lib/spectrumos/colors.conf
PathChanged=/var/lib/spectrumos/current.png
Unit=spectrumos-limine-sync.service

[Install]
WantedBy=multi-user.target
